# rootfs_overlay/etc/init.d/xorg
#!/sbin/openrc-run

description="X11 server (startx)"
command="/usr/bin/xinit"
command_args="/etc/X11/xinit/xinitrc"
command_background=true
pidfile="/var/run/$SVCNAME"

start_pre() {
	ebegin "Generating X11 config"
	X -configure > /dev/null 2>&1
	cp /xorg.conf.new /etc/X11/xorg.conf
	eend $?
}
depend() {
	need localmount
}
