# rootfs_overlay/etc/init.d/xorg
#!/sbin/openrc-run

description="X11 server (startx)"
command="/usr/bin/xinit"
command_args="/etc/X11/xinit/xinitrc"
command_background=true
pidfile="/var/run/$SVCNAME"

depend() {
	need localmount
}
